# deploy/regenerate/sable-kol-regenerate@.service
[Unit]
Description=SableKOL — refresh deliverables for client %i (no SocialData spend)
Documentation=https://github.com/sieggyby/SableKOL/blob/main/docs/sableweb_kol_build_plan.md
After=network-online.target
Wants=network-online.target

# Production install layout (Hetzner CX21):
#   /opt/sable/sable-kol/                — checkout of this repo
#   /opt/sable/sable-kol/.venv/          — Python venv with sable-kol installed
#   /opt/sable/clients/<id>.yaml          — per-client configs (read-only)
#   /opt/sable/outreach/<id>/             — outreach artifacts (writable)
#   /var/log/sable-kol/regenerate-<id>.log — JSON summary per run
#
# Operator install steps:
#   1. Copy this file + the timer to /etc/systemd/system/
#   2. systemctl daemon-reload
#   3. systemctl enable --now sable-kol-regenerate@solstitch.timer
#   4. journalctl -u sable-kol-regenerate@solstitch.service -f to watch
#
# Client_id (the %i specifier) MUST match a YAML at /opt/sable/clients/<id>.yaml.

[Service]
Type=oneshot
User=sable
Group=sable

# Where the CLI lives. Change SABLEKOL_HOME if your install differs.
Environment=SABLEKOL_HOME=/opt/sable/sable-kol
Environment=PYTHONUNBUFFERED=1

# Secrets — ANTHROPIC_API_KEY + SABLE_DATABASE_URL — live in a separate
# root-owned file so they never appear in the unit you might version-
# control. Create it once on the box:
#
#   sudo install -d -m 0750 /etc/sable
#   sudo tee /etc/sable/sable-kol-regenerate.env >/dev/null <<'EOF'
#   ANTHROPIC_API_KEY=sk-ant-...
#   SABLE_DATABASE_URL=postgresql://sable:...@127.0.0.1:5432/sable
#   EOF
#   sudo chmod 0600 /etc/sable/sable-kol-regenerate.env
#   sudo chown root:sable /etc/sable/sable-kol-regenerate.env
#
# Then reload systemd: `sudo systemctl daemon-reload`.
EnvironmentFile=/etc/sable/sable-kol-regenerate.env

WorkingDirectory=/opt/sable/sable-kol

# Single command. --json so systemd journal captures structured output.
ExecStart=/opt/sable/sable-kol/.venv/bin/sable-kol regenerate %i --json

# Hard timeout — Anthropic rate-limits can stretch classify; allocate room
# but cap so a stuck run doesn't hold the next timer.
TimeoutStartSec=2h

# Hardening (Hetzner deploy/DEPLOY.md pattern).
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/sable/outreach /var/log/sable-kol
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictRealtime=true

# Don't restart on failure — the timer fires daily, no need to thrash.
Restart=no

[Install]
WantedBy=multi-user.target
